# files/etc/init.d/overlay_mmc
#!/bin/sh /etc/rc.common
START=15
USE_PROCD=1

start() {
  logger "overlay_mmc: 尝试挂载 TF 卡 ext4 分区作为 overlay"

  # 已挂载则跳过
  mountpoint -q /overlay && {
    logger "overlay_mmc: overlay 已挂载，跳过"
    return 0
  }

  # 找第一个 mmcblk0p* ext4 分区
  MMC_DEV=$(blkid | grep mmcblk0p | grep ext4 | awk -F: '{print $1}' | head -n1)
  if [ -z "$MMC_DEV" ]; then
    logger "overlay_mmc: 未检测到 TF 卡 ext4 分区"
    return 1
  fi

  logger "overlay_mmc: 挂载 $MMC_DEV 到 /mnt"
  mount -t ext4 "$MMC_DEV" /mnt || {
    logger "overlay_mmc: 挂载失败"
    return 1
  }

  # 判断是否已初始化 overlay
  if [ -f /mnt/etc/openwrt_release ]; then
    logger "overlay_mmc: overlay 已初始化，切换挂载"
    umount /overlay || true
    mount --move /mnt /overlay
  else
    logger "overlay_mmc: 初始化 overlay 分区"
    mkdir -p /mnt/upper /mnt/work
    mount -t overlay overlay -o lowerdir=/,upperdir=/mnt/upper,workdir=/mnt/work /mnt
    cp -a /rom/* /mnt/upper/
    sync
    umount /mnt

    # 写入 fstab，保证重启自动挂载
    logger "overlay_mmc: 写入 fstab 以保证自动挂载"
    uci batch <<-EOF
    delete fstab.overlay
    set fstab.overlay=mount
    set fstab.overlay.device=$MMC_DEV
    set fstab.overlay.target=/overlay
    set fstab.overlay.fstype=ext4
    set fstab.overlay.options=defaults
    set fstab.overlay.enabled=1
    commit fstab
EOF
    /etc/init.d/fstab restart

    logger "overlay_mmc: 初始化完成，建议重启"
  fi
}
